# private_dot_config/etc/systemd/system/paccache.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/paccache -r
ExecStart=/usr/bin/paccache -ruk0