# private_dot_config/etc/systemd/system/sshd.service.d/override.conf
[Unit]
Requires=network-online.target
After=
After=sshdgenkeys.service network-online.target