# private_dot_config/systemd/user/delete-pass-history.service
[Unit]
Description=delete password store history

[Service]
Type=oneshot
ExecStart=%h/.local/bin/delete-pass-history
