# private_dot_config/systemd/user/delete-pass-history.timer
[Unit]
Description=delete password store history

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
