# private_dot_config/systemd/user/dnscrypt-adblock.service
[Unit]
Description=Run dnscrypt-adblock

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dnscrypt-adblock
