# private_dot_config/systemd/user/swaywm.target
[Unit]
Description=swaywm systemd user target
Documentation=man:systemd.special(7)
Requires=default.target
After=default.target
AllowIsolate=yes
