# private_dot_config/systemd/user/syncthing-notify.service.tmpl
[Unit]
Description=notify if files are present in syncthing folder

[Service]
Type=oneshot
ExecStart=%h/.local/bin/monitor-dir {{ .syncthing_dir }}
