# private_dot_config/systemd/user/spotifyd.service.d/override.conf
[Unit]
After=
After=sound.target network-online.target swaywm.target

[Install]
WantedBy=
WantedBy=swaywm.target