# private_dot_config/systemd/user/spotifyd.service.d/override.conf
[Unit]
Wants=
After=
Wants=sound.target network-online.service
After=sound.target network-online.service